# dev/build.image/packaging/debuild/scripts/openliberty@.service
[Unit]
Description=Open Liberty server
After=network.target
After=systemd-user-sessions.service
After=network-online.target

[Service]
Type=forking

ExecStart=/etc/init.d/openliberty start %i
ExecStop=/etc/init.d/openliberty stop %i

User=openliberty
Group=openliberty

[Install]
WantedBy=multi-user.target