# dev/com.ibm.ws.security.wim.adapter.ldap_fat.krb5.1/publish/files/ldap-server/supervisord.conf
[inet_http_server]
port=9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9001

[program:slapd]
priority=1
command = /usr/sbin/slapd -h "ldap:/// ldapi:///" -u ldap -g ldap -d 1
startretries = 1
startsecs = 5
redirect_stderr=true
stdout_logfile=/etc/ldapstdout.log
stderr_logfile=/etc/ldapstderr.log

[supervisord]
logfile = /tmp/supervisord_zbx_server.log
loglevel = critical
nodaemon = true
user = root
pidfile = /tmp/supervisord_zbx_server.pid
directory = /tmp
